# roles/monitoring_prometheus/templates/prometheus.service.j2

[Unit]
Description=Prometheus Server
Documentation=https://prometheus.io/docs/introduction/overview/
After=network-online.target

[Service]
User=prometheus
Group=prometheus
Type=simple
Restart=on-failure
RestartSec=5s
ExecStart=/usr/local/bin/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/var/lib/prometheus/data \
  --web.console.templates=/etc/prometheus/consoles \
  --web.console.libraries=/etc/prometheus/console_libraries \
  --web.listen-address=0.0.0.0:{{ prometheus_port | default(9090) }} \
  --web.external-url=http://{{ ansible_default_ipv4.address }}:{{ prometheus_port | default(9090) }}/ \
  --storage.tsdb.retention.time={{ prometheus_retention_time | default('15d') }} \
  --storage.tsdb.retention.size={{ prometheus_retention_size | default('10GB') }} \
  --web.enable-lifecycle \
  --web.enable-admin-api

ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=20s
SendSIGKILL=no

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/prometheus
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target

