# roles/reliability_monitor/templates/reliability-monitor.service.j2

[Unit]
Description=HX Infrastructure Reliability Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/reliability_monitor/health_checker.py --type all --config /etc/reliability_monitor/config.yml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reliability-monitor

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ansible /var/lib/reliability_monitor
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
