# roles/backup_automation/templates/database-backup-metrics.service.j2

[Unit]
Description=Database Backup Metrics Collection Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=backup
Group=backup
ExecStart=/usr/bin/python3 /var/lib/backup/scripts/collect-database-backup-metrics.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ backup_automation.base_directory }} {{ backup_automation.log_directory }} /var/lib/backup
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
MemoryLimit=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
